# systemd/user/wob.socket
[Socket]
ListenFIFO=%t/wob.fifo
SocketMode=0600

[Install]
WantedBy=sockets.target
RequiredBy=wob.service
